# system_files/usr/lib/systemd/system/systemd-suspend.service.d/armada-fake-suspend.conf
# Real S3 hangs on SM8550; redirect suspend requests to fake-suspend.
[Service]
ExecStart=
ExecStart=/usr/libexec/armada/fake-suspend sleep
