# system_files/usr/lib/systemd/system/bootc-generic-growpart.service.d/10-armada-bare-metal.conf
# Armada flashes fixed-size images to bare metal SD cards; grow the btrfs rootfs.

[Unit]
ConditionVirtualization=
