# system_files/usr/lib/systemd/system/armada-perf-paths.service
[Unit]
Description=armada: discover dynamic CPU/GPU/thermal sysfs paths
DefaultDependencies=no
After=local-fs.target
Before=basic.target multi-user.target
RequiresMountsFor=/run

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/armada/perf-paths
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=basic.target
