# system_files/usr/lib/systemd/system/armada-powerkey.service
[Unit]
Description=armada power-button watcher (fake suspend toggle)
After=multi-user.target

[Service]
ExecStart=/usr/libexec/armada/powerkey-watch
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
