# system_files/usr/lib/systemd/system/armada-save-devcoredump@.service
[Unit]
Description=Save kernel devcoredump %i before auto-delete

[Service]
Type=oneshot
TimeoutStartSec=120
ExecStart=/usr/libexec/armada/save-devcoredump %i
