# system_files/usr/lib/systemd/system/armada-steamapps.service
[Unit]
Description=armada: create the Steam library as a btrfs nodatacow subvolume
After=local-fs.target
RequiresMountsFor=/var/home/armada
Before=sddm.service
ConditionPathExists=!/var/home/armada/.local/share/Steam/steamapps

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/armada/setup-steamapps

[Install]
WantedBy=multi-user.target
