# system_files/usr/lib/systemd/system/armada-game-watch.service
[Unit]
Description=armada: apply game power policy while a Steam game runs
After=armada-perf-paths.service
Requires=armada-perf-paths.service

[Service]
ExecStart=/usr/libexec/armada/game-watch
ExecStopPost=/usr/libexec/armada/power-control balanced
Restart=on-failure

[Install]
WantedBy=multi-user.target
